# salt/roots/salt/maxsim/files/supervisor/maxsim.conf
[program:maxsim]
command={{ command }}
directory={{ directory }}
user={{ user }}
autostart=true
autorestart=true
startretries=3
redirect_stderr=True
